# deployment/url-shortener.service
[Unit]
Description=URL Shortener FastAPI App
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/url-shortener
Environment=PATH=/var/www/url-shortener/venv/bin
ExecStart=/var/www/url-shortener/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/url-shortener

[Install]
WantedBy=multi-user.target
